# data/serialosc.service
[Unit]
Description=osc server for monome devices

[Service]
ExecStart=/usr/bin/serialoscd

[Install]
WantedBy=multi-user.target
